# build/config/includes.chroot/etc/skel/.config/systemd/user/hn-wallpaper.service
[Unit]
Description=Rotate HackNow wallpaper

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hn-wallpaper-rotate
