# build/config/includes.chroot/etc/skel/.config/autostart/conky-network.desktop
[Desktop Entry]
Type=Application
Name=Conky Network Intel
Exec=sh -c "sleep 6 && conky -c ~/.config/conky/conky-network.conf"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
