# build/config/includes.chroot/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky System Monitor
Exec=sh -c "sleep 5 && conky -c ~/.config/conky/conky.conf"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
