# build/config/includes.chroot/etc/skel/.config/autostart/hacknow-first-setup.desktop
[Desktop Entry]
Type=Application
Name=HackNow OS First Setup
Exec=hacknow-first-setup
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
